# package/harvester-os/files/etc/systemd/system/getty@tty1.service.d/override.conf
[Service]
EnvironmentFile=-/etc/rancher/installer/env
ExecStartPre=/usr/bin/setterm --msg off
ExecStart=
ExecStart=-/sbin/agetty -n -l /usr/bin/start-harvester-installer.sh %I $TERM
